# sys/bonfire-process.conf
start on filesystem and net-device-up IFACE=lo

stop on shutdown

respawn

script
  . /home/apps/env/bonfire/bin/activate
  exec bonfire process womenwhocode
end script
